# rootfs/etc/systemd/system/sailtrack-timesync.service
[Unit]
Description=SailTrack-TimeSync

[Service]
EnvironmentFile=-/etc/default/sailtrack
Restart=on-failure
ExecStart=/boot/sailtrack/sailtrack-timesync

[Install]
WantedBy=multi-user.target
