# rootfs/etc/systemd/system/sailtrackd-status.service
[Unit]
Description=SailTrackd - Status

[Service]

[Install]
