# rootfs/etc/systemd/system/sailtrackd.service
[Unit]
Description=SailTrackd

[Service]

[Install]
