# rootfs/etc/systemd/system/sailtrack-sailtrackd_status.service
[Unit]
Description=Sailtrack - Sailtrackd Status

[Service]

[Install]
