# rootfs/etc/systemd/system/sailtrack-timesync.service
[Unit]
Description=SailTrack-TimeSync

[Service]
ExecStart=/boot/sailtrack/sailtrack-timesync.py

[Install]
WantedBy=multi-user.target
