# rootfs/etc/systemd/system/sailtrack-status.service
[Unit]
Description=SailTrack-Status

[Service]
Restart=on-failure
ExecStart=/boot/sailtrack/sailtrack-status

[Install]
WantedBy=multi-user.target
