# rootfs/etc/systemd/system/sailtrackd-filter.service
[Unit]
Description=SailTrackd - Filter

[Service]

[Install]
