# rootfs/etc/systemd/system/sailtrackd-imu.service
[Unit]
Description=SailTrackd-IMU

[Service]

[Install]
